# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f

[program:tomcat]
command=/bin/sh -c "source /etc/default/tomcat7 && /usr/share/tomcat7/bin/catalina.sh run"
