# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f

[program:tomcat]
command=/bin/sh -c "source /etc/default/tomcat6 && /usr/share/tomcat6/bin/catalina.sh run"
